# Test/laravel-work.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /Users/Test/Test/artisan queue:work redis --queue=pintu:tag-label
user=www
autostart=true
autorestart=true